# .docker/server/nginx/supervisor/debian_ubuntu/etc/supervisor/conf.d/websockets.conf
[program:websockets]
command=/usr/bin/php /var/www/html/artisan websockets:serve --port=6001
numprocs=1
autostart=true
autorestart=true
user=laravel-echo
